# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=service nginx start
numproc=1
directory=/tmp
stdout_logfile=/var/log/supervisor/nginx.log
autostart=true
user=www-data
stopsignal=KILL
