# etc/supervisor/conf.d/php5-fpm.conf
[program:php5-fpm]
command=service php5-fpm start
numproc=1
directory=/tmp
stdout_logfile=/var/log/supervisor/php5-fpm.log
autostart=true
user=www-data
stopsignal=KILL
